# systemd/LinphoneUI-daemon.service
[Unit]
Description=LinphoneUI Daemon Service
After=dbus.socket pulseaudio.service
Wants=pulseaudio.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/LinphoneUI/scripts/LinphoneUI-daemon.py
Restart=always
RestartSec=5
Environment=PYTHONPATH=/usr/share/LinphoneUI/scripts
Environment=XDG_RUNTIME_DIR=/run/user/%U
WorkingDirectory=/usr/share/LinphoneUI/scripts
User=defaultuser
Group=defaultuser
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target